# infra/systemd/alf-prod.target
[Unit]
Description=Alf Prod Stack
Wants=alf-prod-backend.service

[Install]
WantedBy=default.target
