# infra/systemd/alf-prod-backend.service
[Unit]
Description=Alf Prod — Backend
After=alf-prod-relay.service
Requires=alf-prod-relay.service

[Service]
Type=simple
WorkingDirectory=%h/repos/alf/backend
ExecStart=%h/repos/alf/infra/scripts/run-with-profile.sh %h/repos/alf/backend/node_modules/.bin/tsx src/index.ts
EnvironmentFile=%h/repos/alf/infra/.env.prod
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=alf-prod.target
