# infra/systemd/alf-test.target
[Unit]
Description=Alf Test Stack
Wants=alf-test-relay.service alf-test-backend.service alf-test-frontend.service

[Install]
WantedBy=default.target
