# infra/systemd/alf-prod.target
[Unit]
Description=Alf Prod Stack
Wants=alf-prod-relay.service alf-prod-backend.service alf-prod-frontend.service

[Install]
WantedBy=default.target
