# infra/systemd/alf-dev-frontend.service
[Unit]
Description=Alf Dev — Frontend
After=alf-dev-relay.service

[Service]
Type=simple
WorkingDirectory=%h/repos/alf/frontend
ExecStart=%h/repos/alf/frontend/node_modules/.bin/vite
EnvironmentFile=%h/repos/alf/infra/.env.dev
Environment=VITE_RELAY_URL=ws://localhost:3100/client
# VITE_RELAY_TOKEN is picked up from EnvironmentFile (.env.dev) — vite exposes VITE_* to the browser
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=alf-dev.target
