# infra/systemd/alf-dev-backend.service
[Unit]
Description=Alf Dev — Backend
After=alf-dev-relay.service
Requires=alf-dev-relay.service

[Service]
Type=simple
WorkingDirectory=%h/repos/alf/backend
ExecStart=%h/repos/alf/backend/node_modules/.bin/tsx src/index.ts
EnvironmentFile=%h/repos/alf/infra/.env.dev
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=alf-dev.target
